# deploy/lensing-host.service
[Unit]
Description=Lensing Host Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/lensing
ExecStart=/home/pi/lensing/node_modules/.bin/tsx /home/pi/lensing/packages/cli/src/bin/dev-host.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lensing-host

[Install]
WantedBy=multi-user.target
